# agent/rootfs/etc/s6-overlay/s6-rc.d/svc-desktop/run
#!/command/with-contenv bash

export DISPLAY=:0
export OPENCLAW_USER="${CLAWORC_OPENCLAW_USER:-claworc}"
export OPENCLAW_HOME="${CLAWORC_OPENCLAW_HOME:-/home/claworc}"
export HOME="$OPENCLAW_HOME"
export LANG=en_US.UTF-8
export BROWSER_DATA_DIR="${OPENCLAW_BROWSER_DATA_DIR:-$OPENCLAW_HOME/chrome-data}"
export BROWSER_METRICS_DIR="$BROWSER_DATA_DIR/DeferredBrowserMetrics"

while ! xdpyinfo -display :0 >/dev/null 2>&1; do sleep 0.5; done

s6-setuidgid "$OPENCLAW_USER" openbox &

prepare_browser_state() {
  mkdir -p "$BROWSER_METRICS_DIR"
  chown -R "$OPENCLAW_USER:$OPENCLAW_USER" "$BROWSER_DATA_DIR"
  find "$BROWSER_METRICS_DIR" -mindepth 1 -delete 2>/dev/null || true
}

EXTRA_ARGS=()
if [ -n "$CHROMIUM_USER_AGENT" ]; then
  EXTRA_ARGS+=("--user-agent=$CHROMIUM_USER_AGENT")
fi

if command -v google-chrome-stable >/dev/null 2>&1; then
  BROWSER=google-chrome-stable
  EXTRA_ARGS+=(--no-sandbox)
elif command -v brave-browser >/dev/null 2>&1; then
  BROWSER=brave-browser
  EXTRA_ARGS+=(--no-sandbox)
else
  BROWSER=/usr/bin/chromium
  EXTRA_ARGS+=(--no-sandbox)
fi

while true; do
  prepare_browser_state
  s6-setuidgid "$OPENCLAW_USER" $BROWSER \
    --no-first-run \
    --password-store=basic \
    --simulate-outdated-no-au='Tue, 31 Dec 2099 23:59:59 GMT' \
    --start-maximized \
    --user-data-dir="$BROWSER_DATA_DIR" \
    --remote-debugging-port=9222 \
    --remote-debugging-address=127.0.0.1 \
    --disable-default-apps \
    --disable-breakpad \
    --disable-features=CloseWindowWithLastTab \
    --disable-blink-features=AutomationControlled \
    --disable-infobars \
    --disable-component-update \
    --load-extension=/opt/stealth-extension \
    "${EXTRA_ARGS[@]}" > /dev/null 2>&1
  sleep 1
done
