# agent/rootfs/etc/s6-overlay/s6-rc.d/init-setup/up
/bin/bash -c "
set -e
mkdir -p /var/log/claworc
chmod 755 /var/log/claworc
if [ ! -f /home/claworc/.bashrc ]; then
  cp -a /etc/skel/. /home/claworc/
fi
mkdir -p /home/claworc/chrome-data /home/claworc/.openclaw
chown -R claworc:claworc /home/claworc
"
